# etc/systemd/user/hidpipe.socket
[Unit]
Description=hidpipe daemon socket

[Socket]
ListenStream=%t/hidpipe
NoDelay=true

[Install]
WantedBy=sockets.target
